# root/etc/uci-defaults/30_luci-theme-argon-noodles
#!/bin/sh
uci batch <<-EOF
	set luci.themes.noodlesLight_Purple=/luci-static/argon_noodles
	set luci.themes.noodlesLight=/luci-static/argon_noodles1
	commit luci
EOF
exit 0
