# root/etc/uci-defaults/30_luci-theme-argon-noodles
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon_noodles=/luci-static/argon_noodles
	commit luci
EOF
exit 0
